# Miniconda/pkgs/libmamba-2.0.5-hcd6fe79_1/Library/etc/profile.d/mamba.sh
export MAMBA_ROOT_PREFIX="C:/b/abs_627vsv8bhu/croot/mamba-split_1734469608328/_h_env/Library"
__mamba_setup="$("C:/b/abs_627vsv8bhu/croot/mamba-split_1734469608328/_h_env/Library/bin/mamba" shell hook --shell posix 2> /dev/null)"
if [ $? -eq 0 ]; then
    eval "$__mamba_setup"
else
    alias mamba="C:/b/abs_627vsv8bhu/croot/mamba-split_1734469608328/_h_env/Library/bin/mamba"  # Fallback on help from mamba activate
fi
unset __mamba_setup
